# deploy/systemd/gonba-vk-sync.service
[Unit]
Description=GONBA VK Auto-Sync trigger
After=network-online.target gonba.service
Wants=network-online.target

[Service]
Type=oneshot
User=valstan
EnvironmentFile=-/home/valstan/GONBA/web/.env
ExecStart=/usr/bin/curl --fail --silent --show-error --max-time 600 -X POST \
  -H "Authorization: Bearer ${CRON_SECRET}" \
  -H "Content-Type: application/json" \
  --data "{}" \
  http://127.0.0.1:3000/api/vk-auto-sync/trigger
